# crontab
# min	hour	day	month	weekday	command
*   4	*	*	*	python /srv/miuitask.py